# lava/systemd-system/certbot.timer
[Unit]
Description=Renewal of Certbot certificates

[Timer]
OnCalendar=Daily
#OnCalendar=Weekly
Persistent=true

[Install]
WantedBy=timers.target
