# lava/systemd-system/gnusocial-maintenance.timer
[Unit]
Description=GNU Social maintenance

[Timer]
#OnCalendar=Daily
OnCalendar=Hourly
Persistent=true

[Install]
WantedBy=timers.target
